# reprodebian/cocoon-init/cocoon-cert-refresh.service
[Unit]
Description=Refresh TEE certificate
After=cocoon-ready.target
Requires=cocoon-ready.target

[Service]
Type=oneshot
ExecStart=/usr/bin/gen-cert-synced --tee tee --name /etc/tee/tee --force
