# reprodebian/cocoon-init/cocoon-sev-pki-update.service
[Unit]
Description=SEV PKI update
After=cocoon-ready.target
Requires=cocoon-ready.target

[Service]
Type=oneshot
ExecStart=/usr/bin/sev-pki-update
